# system_files/etc/xdg/autostart/armada-desktop-bootstrap.desktop
[Desktop Entry]
Type=Application
Name=Armada Desktop Bootstrap
Exec=/usr/libexec/armada/desktop-bootstrap
Terminal=false
