# system_files/usr/lib/systemd/system/armada-bootimg-sync.service
[Unit]
Description=Rebuild the ABL boot.img for the next-boot deployment before reboot
ConditionPathExists=/run/ostree-booted
DefaultDependencies=no
RequiresMountsFor=/sysroot /boot /boot/efi
# Reverse stop order: ordered before finalize so our ExecStop runs after it,
# against the finalized next-boot deployment.
After=local-fs.target
Before=ostree-finalize-staged.service final.target
Conflicts=final.target

[Service]
Type=oneshot
RemainAfterExit=yes
# No-op start so the oneshot activates and its ExecStop fires at shutdown.
ExecStart=/bin/true
ExecStop=/usr/libexec/armada/armada-bootimg-update

[Install]
WantedBy=multi-user.target
